# folio-mod-source-record-storage/orig-packaging/folio-mod-source-record-storage.target
[Unit]
Description=FOLIO mod-source-record-storage target

[Install]
WantedBy=multi-user.target
