# folio-mod-patron/orig-packaging/folio-mod-patron.target
[Unit]
Description=FOLIO mod-patron target

[Install]
WantedBy=multi-user.target
