# folio-mod-permissions/orig-packaging/folio-mod-permissions.target
[Unit]
Description=FOLIO mod-permissions target

[Install]
WantedBy=multi-user.target
