# folio-mod-email/orig-packaging/folio-mod-email.target
[Unit]
Description=FOLIO mod-email target

[Install]
WantedBy=multi-user.target
