# folio-mod-event-config/orig-packaging/folio-mod-event-config.target
[Unit]
Description=FOLIO mod-event-config target

[Install]
WantedBy=multi-user.target
