# folio-mod-data-import-converter-storage/orig-packaging/folio-mod-data-import-converter-storage.target
[Unit]
Description=FOLIO mod-data-import-converter-storage target

[Install]
WantedBy=multi-user.target
