# folio-mod-inventory-storage/orig-packaging/folio-mod-inventory-storage.target
[Unit]
Description=FOLIO mod-inventory-storage target

[Install]
WantedBy=multi-user.target
