# folio-mod-erm-usage-harvester/orig-packaging/folio-mod-erm-usage-harvester.target
[Unit]
Description=FOLIO mod-erm-usage-harvester target

[Install]
WantedBy=multi-user.target
