# folio-mod-ncip/orig-packaging/folio-mod-ncip.target
[Unit]
Description=FOLIO mod-ncip target

[Install]
WantedBy=multi-user.target
