# folio-mod-users-bl/orig-packaging/folio-mod-users-bl.target
[Unit]
Description=FOLIO mod-users-bl target

[Install]
WantedBy=multi-user.target
