# folio-mod-gobi/orig-packaging/folio-mod-gobi.target
[Unit]
Description=FOLIO mod-gobi target

[Install]
WantedBy=multi-user.target
