# folio-mod-licenses/orig-packaging/folio-mod-licenses.target
[Unit]
Description=FOLIO mod-licenses target

[Install]
WantedBy=multi-user.target
