# folio-mod-organizations/orig-packaging/folio-mod-organizations.target
[Unit]
Description=FOLIO mod-organizations target

[Install]
WantedBy=multi-user.target
