# folio-mod-test/orig-packaging/folio-mod-test.target
[Unit]
Description=FOLIO mod-test target

[Install]
WantedBy=multi-user.target
