# folio-mod-pubsub/orig-packaging/folio-mod-pubsub.target
[Unit]
Description=FOLIO mod-pubsub target

[Install]
WantedBy=multi-user.target
