# folio-mod-template-engine/orig-packaging/folio-mod-template-engine.target
[Unit]
Description=FOLIO mod-template-engine target

[Install]
WantedBy=multi-user.target
