# folio-mod-rtac/orig-packaging/folio-mod-rtac.target
[Unit]
Description=FOLIO mod-rtac target

[Install]
WantedBy=multi-user.target
