# folio-mod-codex-inventory/orig-packaging/folio-mod-codex-inventory.target
[Unit]
Description=FOLIO mod-codex-inventory target

[Install]
WantedBy=multi-user.target
