# folio-mod-codex-ekb/orig-packaging/folio-mod-codex-ekb.target
[Unit]
Description=FOLIO mod-codex-ekb target

[Install]
WantedBy=multi-user.target
