# folio-mod-invoice-storage/orig-packaging/folio-mod-invoice-storage.target
[Unit]
Description=FOLIO mod-invoice-storage target

[Install]
WantedBy=multi-user.target
