# folio-mod-login-saml/orig-packaging/folio-mod-login-saml.target
[Unit]
Description=FOLIO mod-login-saml target

[Install]
WantedBy=multi-user.target
