# folio-mod-codex-mux/orig-packaging/folio-mod-codex-mux.target
[Unit]
Description=FOLIO mod-codex-mux target

[Install]
WantedBy=multi-user.target
