# folio-mod-notify/orig-packaging/folio-mod-notify.target
[Unit]
Description=FOLIO mod-notify target

[Install]
WantedBy=multi-user.target
