# folio-mod-inventory/orig-packaging/folio-mod-inventory.target
[Unit]
Description=FOLIO mod-inventory target

[Install]
WantedBy=multi-user.target
