# folio-mod-courses/orig-packaging/folio-mod-courses.target
[Unit]
Description=FOLIO mod-courses target

[Install]
WantedBy=multi-user.target
