# folio-mod-kb-ebsco-java/orig-packaging/folio-mod-kb-ebsco-java.target
[Unit]
Description=FOLIO mod-kb-ebsco-java target

[Install]
WantedBy=multi-user.target
