# folio-mod-organizations-storage/orig-packaging/folio-mod-organizations-storage.target
[Unit]
Description=FOLIO mod-organizations-storage target

[Install]
WantedBy=multi-user.target
