# folio-mod-sender/orig-packaging/folio-mod-sender.target
[Unit]
Description=FOLIO mod-sender target

[Install]
WantedBy=multi-user.target
