# folio-mod-circulation/orig-packaging/folio-mod-circulation.target
[Unit]
Description=FOLIO mod-circulation target

[Install]
WantedBy=multi-user.target
