# folio-mod-data-import/orig-packaging/folio-mod-data-import.target
[Unit]
Description=FOLIO mod-data-import target

[Install]
WantedBy=multi-user.target
