# folio-mod-tags/orig-packaging/folio-mod-tags.target
[Unit]
Description=FOLIO mod-tags target

[Install]
WantedBy=multi-user.target
