# folio-mod-feesfines/orig-packaging/folio-mod-feesfines.target
[Unit]
Description=FOLIO mod-feesfines target

[Install]
WantedBy=multi-user.target
