# folio-mod-circulation-storage/orig-packaging/folio-mod-circulation-storage.target
[Unit]
Description=FOLIO mod-circulation-storage target

[Install]
WantedBy=multi-user.target
