# folio-mod-source-record-manager/orig-packaging/folio-mod-source-record-manager.target
[Unit]
Description=FOLIO mod-source-record-manager target

[Install]
WantedBy=multi-user.target
