# folio-mod-login/orig-packaging/folio-mod-login.target
[Unit]
Description=FOLIO mod-login target

[Install]
WantedBy=multi-user.target
