# folio-mod-configuration/orig-packaging/folio-mod-configuration.target
[Unit]
Description=FOLIO mod-configuration target

[Install]
WantedBy=multi-user.target
