# folio-mod-finance-storage/orig-packaging/folio-mod-finance-storage.target
[Unit]
Description=FOLIO mod-finance-storage target

[Install]
WantedBy=multi-user.target
