# folio-mod-patron-blocks/orig-packaging/folio-mod-patron-blocks.target
[Unit]
Description=FOLIO mod-patron-blocks target

[Install]
WantedBy=multi-user.target
