# folio-mod-finance/orig-packaging/folio-mod-finance.target
[Unit]
Description=FOLIO mod-finance target

[Install]
WantedBy=multi-user.target
