# folio-mod-orders-storage/orig-packaging/folio-mod-orders-storage.target
[Unit]
Description=FOLIO mod-orders-storage target

[Install]
WantedBy=multi-user.target
