# folio-mod-password-validator/orig-packaging/folio-mod-password-validator.target
[Unit]
Description=FOLIO mod-password-validator target

[Install]
WantedBy=multi-user.target
