# folio-mod-calendar/orig-packaging/folio-mod-calendar.target
[Unit]
Description=FOLIO mod-calendar target

[Install]
WantedBy=multi-user.target
