# folio-mod-erm-usage/orig-packaging/folio-mod-erm-usage.target
[Unit]
Description=FOLIO mod-erm-usage target

[Install]
WantedBy=multi-user.target
