# folio-mod-users/orig-packaging/folio-mod-users.target
[Unit]
Description=FOLIO mod-users target

[Install]
WantedBy=multi-user.target
