# folio-mod-quick-marc/orig-packaging/folio-mod-quick-marc.target
[Unit]
Description=FOLIO mod-quick-marc target

[Install]
WantedBy=multi-user.target
