# folio-mod-notes/orig-packaging/folio-mod-notes.target
[Unit]
Description=FOLIO mod-notes target

[Install]
WantedBy=multi-user.target
