# folio-mod-user-import/orig-packaging/folio-mod-user-import.target
[Unit]
Description=FOLIO mod-user-import target

[Install]
WantedBy=multi-user.target
